# slh/scripts/slh.timer
[Unit]
Description=Stmarks library Helper server runner

[Timer]
OnCalendar=8:00
Persistent=True

[Install]
WantedBy=multi-user.target
