# timer/crontab
00 00 * * * node /srv/lib/index.js
00 08 * * * node /srv/lib/index.js
