# timer/crontab
00 00 00,08 * * * node /srv/src/index.js