# timer/crontab
00 00 * * * node /srv/src/index.js
00 08 * * * node /srv/src/index.js
